# docker/single-container/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/run/supervisord.pid
user=root

[program:mariadb]
priority=10
command=/usr/sbin/mariadbd --user=mysql --datadir=/var/lib/mysql --socket=/run/mysqld/mysqld.sock --pid-file=/run/mysqld/mysqld.pid --bind-address=127.0.0.1 --port=3306 --character-set-server=utf8mb4 --collation-server=utf8mb4_unicode_ci
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/data/logs/mariadb-stdout.log
stderr_logfile=/data/logs/mariadb-stderr.log
stopasgroup=true
killasgroup=true

[program:backend]
priority=20
command=/usr/local/bin/start-backend.sh
autostart=true
autorestart=true
startsecs=10
startretries=10
stdout_logfile=/data/logs/backend-stdout.log
stderr_logfile=/data/logs/backend-stderr.log
stopasgroup=true
killasgroup=true

[program:nginx]
priority=30
command=/usr/local/bin/start-nginx.sh
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/data/logs/nginx-stdout.log
stderr_logfile=/data/logs/nginx-stderr.log
stopasgroup=true
killasgroup=true
